# tester/events-launcher.service
[Unit]
Discription=test launcher

[Service]

ExecStart=/bin/bash -c '{PWD}/launcher/launcher >> {PWD}/launcher.log &'
ExecStop=/bin/kill -s QUIT $MAINPID

Restart=always
RestartSec=20
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
